# deploy/supervisor.conf
; 画己职测 — Supervisor 进程管理配置
;
; 品牌：画己职测
; 项目：career_test
;
; 管理进程：
;   1. career_test_gunicorn      — Gunicorn Web 服务
;   2. career_test_celery_worker — Celery 异步任务 Worker
;   3. career_test_celery_beat   — Celery 定时任务调度器
;
; 使用方式：
;   sudo cp deploy/supervisor.conf /etc/supervisor/conf.d/career_test.conf
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start career_test_gunicorn career_test_celery_worker career_test_celery_beat
;
; 变量说明（请根据实际部署路径修改）：
;   %(here)s 指向本配置文件所在目录，需替换为项目实际路径

; ===========================================================================
; Gunicorn — Web 应用服务（Nginx → Gunicorn → Django）
; ===========================================================================
[program:career_test_gunicorn]
command=/var/www/career_test/venv/bin/gunicorn career_test.wsgi:application --config /var/www/career_test/deploy/gunicorn.conf.py
directory=/var/www/career_test
user=www-data
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
killasgroup=true
stopasgroup=true
stderr_logfile=/var/log/career_test/gunicorn_err.log
stdout_logfile=/var/log/career_test/gunicorn_out.log
environment=DJANGO_SETTINGS_MODULE="career_test.settings.prod"
priority=10

; ===========================================================================
; Celery Worker — 异步任务执行（订单过期、数据清理、报告生成等）
; ===========================================================================
[program:career_test_celery_worker]
command=/var/www/career_test/venv/bin/celery -A career_test worker --loglevel=info --concurrency=4
directory=/var/www/career_test
user=www-data
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=60
killasgroup=true
stopasgroup=true
stderr_logfile=/var/log/career_test/celery_worker_err.log
stdout_logfile=/var/log/career_test/celery_worker_out.log
environment=DJANGO_SETTINGS_MODULE="career_test.settings.prod"
priority=20

; ===========================================================================
; Celery Beat — 定时任务调度器
; ===========================================================================
[program:career_test_celery_beat]
command=/var/www/career_test/venv/bin/celery -A career_test beat --loglevel=info --scheduler django_celery_beat.schedulers:DatabaseScheduler
directory=/var/www/career_test
user=www-data
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
killasgroup=true
stopasgroup=true
stderr_logfile=/var/log/career_test/celery_beat_err.log
stdout_logfile=/var/log/career_test/celery_beat_out.log
environment=DJANGO_SETTINGS_MODULE="career_test.settings.prod"
priority=30

; ===========================================================================
; 日志目录（Supervisor 启动前需确保 /var/log/career_test 存在）
; ===========================================================================
[group:career_test]
programs=career_test_gunicorn,career_test_celery_worker,career_test_celery_beat
